# server-health-check.service
[Unit]
Description=Server Health Check

[Service]
ExecStart=/bin/ping -c 3 example.com
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
